# audioforge.service
[Unit]
Description=audioforge — AI Music Generator
After=network.target

[Service]
Type=simple
User=marc
WorkingDirectory=/home/marc/music-generator
ExecStart=/home/marc/music-generator/run.sh 8989
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
